# dot_config/opencode/profiles/web/opencode-serve.service
# OpenCode web server — systemd service unit
#
# This file is the canonical definition for the opencode-serve service. It runs
# opencode in "web" mode (API + Web UI) locked to the "web" profile so that the
# Google-focused provider config (opencode-google.json) and the companion env
# file (service.env) in this directory are used at startup.
#
# Installation (as root or with sudo):
#   cp opencode-serve.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now opencode-serve
#
# The service expects:
#   /home/ubuntu/.config/opencode/profiles/web/service.env     — API keys & credentials
#   /home/ubuntu/.config/opencode/profiles/web/opencode.json   — profile config
#   /home/ubuntu/.opencode/bin/opencode                        — opencode binary
#
# Logs: journalctl -u opencode-serve -f

[Unit]
Description=OpenCode web server (API + Web UI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu

# Environment variables (API keys, OAuth secrets) — keep this alongside the
# profile configs so the web profile is self-contained.
EnvironmentFile=/home/ubuntu/.config/opencode/profiles/web/service.env

# Tells opencode to read its config from the "web" profile directory instead of
# the default ~/.config/opencode/opencode.json.
Environment=OPENCODE_CONFIG_DIR=/home/ubuntu/.config/opencode/profiles/web

ExecStart=/home/ubuntu/.opencode/bin/opencode web --port 38889 --hostname 0.0.0.0
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
