# dot_config/opencode/profiles/web/opencode-serve.service
# OpenCode web server — systemd service unit
#
# This file is the canonical definition for the opencode-serve service. It runs
# opencode in "web" mode (API + Web UI) using the global config at
# ~/.config/opencode/opencode.json.
#
# Installation (as root or with sudo):
#   cp opencode-serve.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now opencode-serve
#
# The service expects:
#   ~/.config/opencode/profiles/web/service.env     — API keys & credentials
#   ~/.config/opencode/opencode.json                — global config
#   ~/.opencode/bin/opencode                        — opencode binary
#
# Logs: journalctl -u opencode-serve -f

[Unit]
Description=OpenCode web server (API + Web UI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=%h

# Environment variables (API keys, OAuth secrets)
EnvironmentFile=%h/.config/opencode/profiles/web/service.env

# Uses global config at ~/.config/opencode/opencode.json (default)
# No OPENCODE_CONFIG_DIR needed

ExecStart=%h/.opencode/bin/opencode web --port 38889 --hostname 0.0.0.0
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
